# conf.d/mudlet.conf
[program:mudlet]
command=/usr/games/mudlet --appimage-extract-and-run
autorestart=true
