# conf.d/mudlet.conf
[program:mudlet]
command=/usr/games/mudlet
autorestart=true
