# conf.d/mudlet.conf
[program:mudlet]
command=/usr/games/mudlet --appimage-extract-and-run
startsecs = 0
autostart=true
startretries=10
autorestart=true
