# conf.d/tigervnc.conf
[program:tigervnc]
command=Xtigervnc :0 -rfbport 5900 -rfbauth /mudlet/.vnc/passwd -geometry %(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)s -depth 24 -localhost no
autorestart=true
priority=10
